# kerio-connect-reindex.service
[Unit]
Description=Kerio Connect Index Reset Service
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/kerio-connect-reindex