# debian/prerm
#!/bin/sh
set -e

# Only run on removal and if systemd is actively running
if [ "$1" = "remove" ] || [ "$1" = "deconfigure" ]; then
    if [ -d /run/systemd/system ]; then
        systemctl disable --now kerio-connect-reindex.timer || true
    fi
fi

exit 0